# pi/brew-ble/systemd/brew-ble-scanner.service
[Unit]
Description=Brew BLE Pill Scanner
After=bluetooth.target
Wants=bluetooth.target

[Service]
Type=simple
WorkingDirectory=/opt/brew-ble
ExecStart=/opt/brew-ble/venv/bin/python /opt/brew-ble/ble_scanner.py
Restart=always
RestartSec=5
# BLE scanning needs CAP_NET_RAW + CAP_NET_ADMIN
AmbientCapabilities=CAP_NET_RAW CAP_NET_ADMIN
CapabilityBoundingSet=CAP_NET_RAW CAP_NET_ADMIN

[Install]
WantedBy=multi-user.target